# tools/nocloud-factory-install/factory-install/systemd/factory-install-config.path
[Unit]
Description=Factory Installation Config Stage
ConditionPathExists=/var/lib/factory-install/enabled
ConditionPathExists=!/var/lib/factory-install/state/config

[Path]
PathModified=/var/lib/factory-install/stage/config

[Install]
WantedBy=factory-install.target
