# tools/nocloud-factory-install/factory-install/systemd/factory-install-setup.path
[Unit]
Description=Factory Installation System Setup Stage
ConditionPathExists=/var/lib/factory-install/enabled
ConditionPathExists=!/var/lib/factory-install/state/setup

[Path]
PathExists=/var/lib/factory-install/stage/setup

[Install]
WantedBy=factory-install.target
