# tools/nocloud-factory-install/factory-install/systemd/factory-install-tests.path
[Unit]
Description=Factory Installation System Tests Stage
ConditionPathExists=/var/lib/factory-install/enabled
ConditionPathExists=!/var/lib/factory-install/state/tests

[Path]
PathModified=/var/lib/factory-install/stage/tests

[Install]
WantedBy=factory-install.target
