# tools/nocloud-factory-install/factory-install/systemd/factory-install-backup.path
[Unit]
Description=Factory Installation Backup Stage
ConditionPathExists=/var/lib/factory-install/enabled
ConditionPathExists=!/var/lib/factory-install/state/backup

[Path]
PathModified=/var/lib/factory-install/stage/backup

[Install]
WantedBy=factory-install.target
