# tools/nocloud-factory-install/factory-install/systemd/factory-install-bootstrap.path
[Unit]
Description=Factory Installation Bootstrap Stage
ConditionPathExists=/var/lib/factory-install/enabled
ConditionPathExists=!/var/lib/factory-install/state/bootstrap

[Path]
PathExists=/var/lib/factory-install/stage/bootstrap

[Install]
WantedBy=factory-install.target
