# tools/nocloud-factory-install/factory-install/systemd/factory-install.target
# Target enabled by cloud-init
[Unit]
Description=Factory Install Target
After=cloud-init.target network-online.target

[Install]
WantedBy=multi-user.target
